# services/rpi-trader-scheduler.service
[Unit]
Description=RPI Trader Scheduler
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/rpi-trader
Environment=PATH=/home/pi/miniconda3/envs/env-scheduler/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/pi/miniconda3/envs/env-scheduler/bin/python -m apps.scheduler.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpi-trader-scheduler

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/rpi-trader

# Environment file
EnvironmentFile=/home/pi/rpi-trader/.env

[Install]
WantedBy=multi-user.target

