# services/rpi-trader-execution-worker.service
[Unit]
Description=RPI Trader Execution Worker
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=andrepi
Group=andrepi
WorkingDirectory=/home/andrepi/rpi-trader
Environment=PATH=/home/andrepi/miniconda3/envs/env-execution/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/andrepi/miniconda3/envs/env-execution/bin/python -m apps.execution_worker.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpi-trader-execution-worker

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/andrepi/rpi-trader

# Environment file
EnvironmentFile=/home/andrepi/rpi-trader/.env

[Install]
WantedBy=multi-user.target

